# deploy/cch-redis-session-puller.service.example
[Unit]
Description=CCH Redis Session Puller
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/cch-redis-session-puller.env
WorkingDirectory=/path/to/cch-redis-session-puller
ExecStart=/path/to/cch-redis-session-puller/deploy/run-puller.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
